# env/esrf/inittab
::sysinit:/etc/rc.d/startup
console::respawn:/sbin/getty -L 9600 console
::shutdown:/etc/rc.d/shutdown
::ctrlaltdel:/sbin/reboot
